# scripts/services/yektayar-backend.service
[Unit]
Description=YektaYar Backend API Server
Documentation=https://github.com/atomicdeploy/yektayar
After=network.target postgresql.service

[Service]
Type=simple
User=__USER__
Group=__GROUP__
WorkingDirectory=__PROJECT_PATH__/packages/backend
Environment="NODE_ENV=production"
Environment="PATH=__BUN_PATH__:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=__PROJECT_PATH__/.env
ExecStart=__BUN_PATH__/bun run src/index.ts
Restart=always
RestartSec=10
StandardOutput=append:/var/log/yektayar/backend.log
StandardError=append:/var/log/yektayar/backend-error.log

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Resource limits
LimitNOFILE=65536
TasksMax=512

[Install]
WantedBy=multi-user.target
